# deployment/price-collector-shadow-signal.service
[Unit]
Description=Experimental Chainlink catch-up shadow-signal worker
After=network-online.target postgresql.service redis-server.service price-collector-binance-futures.service price-collector-polymarket-chainlink.service
Wants=network-online.target redis-server.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=pricecollector
Group=pricecollector
WorkingDirectory=/opt/price-collector
EnvironmentFile=/etc/price-collector/shadow-signal.env
ExecStart=/opt/price-collector/.venv/bin/python -m price_collector.shadow_signal_collector

Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
